# edison_tools/extboard_detect.service
[Unit]
Description=MC_EXTBOARD_DETECT
[Service]
Type=forking
ExecStart=/etc/extboardcheck.sh
[Install]
WantedBy=multi-user.target
